# tools/claude-gateway.service
[Unit]
Description=Claude Code — Radio Gateway tmux session
After=network-online.target radio-gateway.service
Wants=network-online.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Downloads/radio-gateway
Environment=HOME=/home/user
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/bin/bash /home/user/Downloads/radio-gateway/tools/claude-gateway.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
